# apps/agent/systemd/adminhelper-agent.service
[Unit]
Description=AdminHelper Agent — FRPC Sync + Monitoring
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/adminhelper-agent run --once
# < OnUnitActiveSec=300 (timer), but enough for the serial collector timeouts (Docker,
# smartctl 30s/disk) plus the monitor-push retry — 60s let a single dying disk + a push
# retry get SIGTERM'd, losing every report (4.13). Keep in sync with the dynamic unit in
# internal/service/service_linux.go.
TimeoutStartSec=270
