# packaging/systemd/srm-frpc-sync.timer
[Unit]
Description=SRM FRP Config Sync Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
