# packaging/systemd/srm-frpc-sync.service
[Unit]
Description=SRM FRP Config Sync

[Service]
Type=oneshot
ExecStart=/usr/local/bin/srm-frpc-sync
